# recipes-angstrom/angstrom/cpufreq-tweaks/cpu-ondemand.service
[Unit]
Description=Switch to ondemand CPU-freq governor

[Service]
Type=oneshot
ExecStart=/usr/bin/cpufreq-set -g ondemand
