# recipes-angstrom/angstrom/cpufreq-tweaks/cpu-schedutil.service
[Unit]
Description=Switch to schedutil CPU-freq governor

[Service]
Type=oneshot
ExecStart=/usr/bin/cpufreq-set -g schedutil
